# deploy/university-desktop/launchd/com.omeia.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.omeia.api</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_VENV/bin/uvicorn</string>
    <string>app_skeleton.api.main:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_REPO_ROOT</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>StandardOutPath</key>
  <string>REPLACE_WITH_REPO_ROOT/logs/omeia-api.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_REPO_ROOT/logs/omeia-api.err.log</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
