# deploy/university-desktop/omeia-quality-eval.timer
[Unit]
Description=Weekly OMEIA continuous quality evaluation

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=1800
Unit=omeia-quality-eval.service

[Install]
WantedBy=timers.target
