# deploy/university-desktop/omeia-processor.service
[Unit]
Description=OMEIA autonomous processor (vault, digitalize, Supabase sync)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=omeia
Group=omeia
WorkingDirectory=/opt/omeia/digital-notepad
EnvironmentFile=/opt/omeia/deploy/university-desktop/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/omeia/venv/bin/python scripts/ops/autonomous_processor.py --daemon
Restart=on-failure
RestartSec=30
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
